# Buraindo IOT/autostart/autorunmyfile2.desktop
[Desktop Entry]
Exec= python3 /home/jayvie2007/Desktop/Capstone/relayHumidity.py