# Buraindo IOT/autostart/autorunmyfile3.desktop
[Desktop Entry]
Exec= sudo python3 /home/jayvie2007/Desktop/Capstone/ledControl.py
